# backends/apt/pk-debconf-helper.socket.in
[Unit]
Description=debconf communication socket

[Socket]
ListenStream=%t/pk-debconf-socket

[Install]
WantedBy=sockets.target
