# backends/aptcc/pk-debconf-helper.service.in
[Unit]
Description=debconf communication service
Requires=pk-debconf-helper.socket
RefuseManualStart=true

[Service]
ExecStart=@libexecdir@/pk-debconf-helper
