# contrib/systemd-updates/packagekit-offline-update.service.in
[Unit]
Description=Updates the operating system whilst offline
OnFailure=reboot.target

[Install]
WantedBy=system-update.target

[Service]
ExecStart=@libexecdir@/pk-offline-update
