# data/packagekit-offline-update.service.in
[Unit]
Description=Updates the operating system whilst offline
OnFailure=reboot.target

[Service]
ExecStart=@libexecdir@/pk-offline-update
