# data/packagekit.service.in
[Unit]
Description=PackageKit Daemon

[Service]
Type=dbus
BusName=org.freedesktop.PackageKit
User=@PACKAGEKIT_USER@
ExecStart=@libexecdir@/packagekitd
