# data/packagekit-offline-update.service.in
[Unit]
Description=Update the operating system whilst offline

DefaultDependencies=no
Requires=sysinit.target dbus.socket
After=sysinit.target dbus.socket systemd-journald.socket
Before=shutdown.target system-update.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/pk-offline-update

FailureAction=reboot
